# config_examples/xanmel.conf
description "upstart configuration for Xanmel Chat bot"

start on net-device-up
stop on shutdown

respawn

console log

setuid xanmel

exec /usr/local/bin/xanmel --config /etc/xanmel.yaml
